# deploy/fold-api.service
# deploy/fold-api.service
[Unit]
Description=fold@Scripps API (uvicorn)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=fold
WorkingDirectory=/opt/fold-at-scripps
EnvironmentFile=/etc/fold/fold.env
ExecStartPre=/usr/bin/env uv run alembic upgrade head
ExecStart=/usr/bin/env uv run uvicorn fold_at_scripps.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
